# services/geoclue-agent.service
[Unit]
Description=Geoclue Agent

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
